# cron.d/pylogan
#Update log processing
*/5 *   * * *   root    /usr/bin/python3 /opt/pylogan/bin/pylogan.py >> /var/log/pylogan.log
