# cron.d/pylogan
#Update log processing
*/5 *   * * *   root    /usr/bin/python3 /usr/local/bin/pylogan_bin.py
